# docker/build/docker-app/crontab
0 6,12,18,0 * * * /var/www/bin/console app:cron-search