# docker/build/docker-app/crontab
* * * * * /var/www/bin/console messenger:consume scheduler_default --limit=1
