# docker/build/docker-app/crontab
0 18 * * * /var/www/bin/console app:cron-search
